# deploy/systemd/worker.service
[Unit]
Description=Trading AI worker compatibility alias
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/trading-ai
Environment="PYTHONUNBUFFERED=1"
ExecStart=/root/trading-ai/.venv/bin/python -m app.workers.research_worker
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
